# docker/etc/service/nginx/run
#!/bin/sh
exec 2>&1
nginx -g 'daemon off;'